# bin/fireplace-server.service
[Unit]
Description=Fireplace Server

[Service]
User=fireplace
Group=fireplace
WorkingDirectory=/opt/fireplace
ExecStart=/opt/fireplace/fireplace-server -host=127.0.0.1:8999
Restart=always

[Install]
WantedBy=multi-user.target